# etc/supervisor/conf.d/sup_rst.conf
[program:mbusd]
command=mbusd
user=root
autostart=true
autorestart=true
startsecs=10
startretries=30

[program:icmpd]
command=icmpd
user=root
autostart=true
autorestart=true
