# vm/supervisor/vm.conf
[program:vm]
command=/opt/vm/run.sh
directory=/opt/vm/
user=victoriametrics
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/vm.log
stdout_logfile_maxbytes=1MB
redirect_stderr=true